# ansible-role-postgresql/templates/postgresql.service.j2
.include /lib/systemd/system/postgresql.service
[Service]

# Location of database directory
Environment=PGDATA={{ postgres_new_data_dir }}
